# docs/launchd/cfo-helper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd template to keep the cfo-helper dashboard running on macOS.

  To install (run once):

    cp docs/launchd/cfo-helper.plist ~/Library/LaunchAgents/com.cfohelper.dashboard.plist
    sed -i '' "s|REPLACE_REPO|$(pwd)|g" ~/Library/LaunchAgents/com.cfohelper.dashboard.plist
    launchctl load ~/Library/LaunchAgents/com.cfohelper.dashboard.plist

  To uninstall:

    launchctl unload ~/Library/LaunchAgents/com.cfohelper.dashboard.plist
    rm ~/Library/LaunchAgents/com.cfohelper.dashboard.plist

  After install, the dashboard auto-starts on login and restarts if it crashes.
  Logs land in REPLACE_REPO/logs/launchd-*.log.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.cfohelper.dashboard</string>

  <key>ProgramArguments</key>
  <array>
    <string>REPLACE_REPO/.venv/bin/uvicorn</string>
    <string>web.main:app</string>
    <string>--port</string>
    <string>8765</string>
    <string>--log-level</string>
    <string>info</string>
  </array>

  <key>WorkingDirectory</key>
  <string>REPLACE_REPO</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>REPLACE_REPO/logs/launchd-stdout.log</string>

  <key>StandardErrorPath</key>
  <string>REPLACE_REPO/logs/launchd-stderr.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
  </dict>
</dict>
</plist>
